# deus-enhanced.service
[Unit]
Description=Deus Ex Machina Enhanced System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/deus-ex-machina
ExecStart=/usr/bin/python3 /opt/deus-ex-machina/enhanced/integration.py monitor --interval 300
Restart=on-failure
RestartSec=5s
# Environment variables for AI provider API keys
Environment="GEMINI_API_KEY="
Environment="ANTHROPIC_API_KEY="
Environment="OPENAI_API_KEY="
# Other environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="DEUS_INSTALL_DIR=/opt/deus-ex-machina"
Environment="DEUS_LOG_DIR=/var/log/deus-ex-machina"

[Install]
WantedBy=multi-user.target
